# mkosi.profiles/kopia-filesystem-backup/mkosi.extra/etc/systemd/system/kopia-filesystem-backup.timer
[Unit]
Description=Hourly Kopia filesystem backup timer

[Timer]
OnCalendar=hourly
# No catch-up: a missed run (drive detached, machine off) is not replayed on
# boot — the next hourly tick is soon enough.
Persistent=false

[Install]
WantedBy=timers.target
