# mkosi.profiles/s3-unencrypted-backup/mkosi.extra/etc/systemd/system/s3-backup.timer
[Unit]
Description=Hourly S3 backup of configured files
Documentation=man:mc(1)

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target