# mkosi.profiles/cloudflare-tunnel/mkosi.extra/etc/systemd/system/cloudflared-watchdog.timer
[Unit]
Description=Periodic cloudflared tunnel health probe

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
AccuracySec=15s
RandomizedDelaySec=30s
Unit=cloudflared-watchdog.service

[Install]
WantedBy=timers.target
