# mkosi.extra/etc/systemd/system/tailscaled.service.d/10-restart.conf
# Tightens the restart policy for the upstream tailscaled.service.
# The upstream unit already has Restart=on-failure; we make the cadence
# explicit and add a start rate-limiter so a hard-looping failure does
# not thrash.

[Service]
Restart=on-failure
RestartSec=10s

# If tailscaled fails 10 times in 10 minutes, stop trying until we are
# reset (systemctl reset-failed tailscaled). The watchdog timer will
# pick it up.
[Unit]
StartLimitIntervalSec=600
StartLimitBurst=10
