# mkosi.extra/etc/systemd/system/ab-dns-ensure.service
[Unit]
Description=Ensure DNS resolver handoff is wired for NetworkManager and Tailscale
Documentation=man:systemd-resolved.service(8) man:resolv.conf(5)
DefaultDependencies=no
After=local-fs.target systemd-resolved.service
Before=network-pre.target NetworkManager.service tailscaled.service tailscale-up.service
Wants=systemd-resolved.service network-pre.target

[Service]
Type=oneshot
# On the sysinit critical path (Before=network-pre.target). Bound it so a stuck
# resolvectl/restart can never hang early boot; the script's own work is <5s.
TimeoutStartSec=30s
ExecStart=/usr/local/libexec/ab-dns-ensure

[Install]
WantedBy=sysinit.target
