# mkosi.extra/etc/systemd/system/ab-dns-ensure.timer
[Unit]
Description=Periodically verify DNS resolver handoff

[Timer]
OnBootSec=30s
OnUnitActiveSec=1min
AccuracySec=30s
Unit=ab-dns-ensure.service

[Install]
WantedBy=timers.target
