# mkosi.extra/etc/systemd/system/ab-grow-root.service
[Unit]
Description=Grow root filesystem to fill its A/B slot
Documentation=file:///usr/local/sbin/ab-grow-root
# Root is already mounted by the time local-fs.target is reached; grow it
# before multi-user.target so logins / apt / services see the full slot.
After=local-fs.target
Before=multi-user.target
ConditionPathExists=/usr/local/sbin/ab-grow-root

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/ab-grow-root
# Ordered Before=multi-user.target, so a wedged growpart/resize2fs blocks the
# login prompt. 5min covers a legitimate large-fs resize on slow media; past
# that, treat it as hung and let the boot proceed rather than wait forever.
TimeoutStartSec=300s

[Install]
WantedBy=multi-user.target
