# mkosi.extra/etc/systemd/system/ab-monitor.service
[Unit]
Description=AB reliability monitor (health checks + alerts)
Documentation=file:///usr/local/libexec/ab-monitor/check.sh
After=network-online.target
Wants=network-online.target
# Run after the main services we are monitoring have had a chance to
# start. We also gate on our own startup grace in check.sh, so this
# is belt-and-suspenders.
After=tailscaled.service cloudflared.service ssh.service
StartLimitIntervalSec=600
StartLimitBurst=20

[Service]
Type=oneshot

# Decrypt whichever outbound-notification credentials are actually
# present in /etc/credstore.encrypted/ into a private runtime
# directory and expose it to check.sh / notify.sh / notifiers/*.sh
# via $CREDENTIALS_DIRECTORY (the same env var systemd's native
# LoadCredentialEncrypted= would set). This matters because:
#   * `LoadCredentialEncrypted=ID:/path` HARD-FAILS the unit if the
#     source path is missing — there is no per-entry "optional" syntax.
#     A clean image (after clean.sh, or before package-alert-credentials
#     has been run) has zero credstore blobs, so the unit would fail
#     to start on every fresh install.
#   * The notifiers under /usr/local/libexec/ab-monitor/notifiers/
#     already tolerate a missing file inside $CREDENTIALS_DIRECTORY,
#     so populating only what exists is exactly what they expect.
RuntimeDirectory=ab-monitor
RuntimeDirectoryMode=0700
Environment=CREDENTIALS_DIRECTORY=%t/ab-monitor
ExecStartPre=/usr/local/libexec/ab-monitor/load-creds-if-present

EnvironmentFile=-/etc/default/ab-monitor
ExecStart=/usr/local/libexec/ab-monitor/check.sh

# We never want the monitor itself to wedge the box.
TimeoutStartSec=120s

# Sandboxing. The monitor needs to read various /run and /var/lib
# paths, run systemctl, call bootctl, hit localhost HTTP, hit the
# public internet. It does NOT need to write to /etc or /usr.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
# StateDirectory= creates /var/lib/ab-monitor before the namespace is
# set up. Without it the unit fails at step NAMESPACE on a fresh image
# (the path simply does not exist yet). RuntimeDirectory=ab-monitor
# above already covers /run/ab-monitor.
StateDirectory=ab-monitor
StateDirectoryMode=0700
ReadWritePaths=/run/log
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
