# mkosi.profiles/cloud-backup/mkosi.extra/etc/systemd/system/cloud-backup.timer
[Unit]
Description=Hourly cloud backup timer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target