# mkosi.profiles/healthchecksio/mkosi.extra/etc/systemd/system/ab-heartbeat.service
[Unit]
Description=Heartbeat to healthchecks.io (dead-man's switch)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=600
StartLimitBurst=20

[Service]
Type=oneshot
# Match ab-monitor.service: tolerate the credstore blob being absent
# on a clean image. heartbeat.sh treats a missing
# $CREDENTIALS_DIRECTORY/healthchecks-ping-url as "no heartbeat
# configured" rather than an error.
RuntimeDirectory=ab-monitor
RuntimeDirectoryMode=0700
Environment=CREDENTIALS_DIRECTORY=%t/ab-monitor
ExecStartPre=/usr/local/libexec/ab-monitor/load-creds-if-present
EnvironmentFile=-/etc/default/ab-monitor
ExecStart=/usr/local/libexec/ab-monitor/heartbeat.sh

TimeoutStartSec=60s

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
RestrictSUIDSGID=yes
