# mkosi.profiles/thinkpad-g13/mkosi.extra/etc/systemd/system/gnome-touchpad-fix.service
[Unit]
Description=Apply GNOME touch‑pad preferences (natural scroll)
After=display-manager.service
Wants=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/bin/gsettings set org.gnome.desktop.peripherals.touchpad natural-scroll true

[Install]
WantedBy=display-manager.service
