# mkosi.profiles/screen-f8-recovery-key/mkosi.extra/etc/systemd/system/triggerhappy.service.d/override.conf
[Service]
# Clear the existing user declaration, then set to root
User=
User=root
