# mkosi.profiles/watchdog-battery/mkosi.extra/etc/systemd/system/watchdog-battery.service
[Unit]
Description=Battery‑watchdog – suspend when charge ≤ 5 %
After=systemd-modules-load.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/watchdog-battery.sh
