# mkosi.profiles/home-server-backup/mkosi.extra/etc/systemd/system/home-server-backup.timer
[Unit]
Description=Hourly home server backup timer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target