# mkosi.profiles/flatpak-support/mkosi.extra/etc/systemd/system/flatpak-flathub-init.service
[Unit]
Description=Add Flathub remote (first-boot, idempotent)
Documentation=https://flatpak.org/setup/Debian
ConditionPathExists=!/var/lib/flatpak/.ab-flathub-added
# Flatpak's apt package creates /var/lib/flatpak/ as part of its
# postinst, so the path must already exist for ExecStartPost's touch
# to succeed.
ConditionPathIsDirectory=/var/lib/flatpak

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/flatpak-flathub-init
